# root/etc/services.d/btsync/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/bin/rslsync --webui.listen 0.0.0.0:$PORT --config /config/rslsync.conf --nodaemon
