# root/etc/services.d/btsync/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/bin/btsync \
		--webui.listen 0.0.0.0:$PORT \
		--storage /config \
		--config /config/btsync.conf