# root/etc/services.d/btsync/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /rslsync/rslsync --config /rslsync/rslsync.conf &> /dev/null
