# root/etc/services.d/btsync/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid root /rslsync/rslsync --config /rslsync/rslsync.conf
